# server.conf
[program:forumserver]
command=/usr/local/bin/gunicorn wsgi --bind localhost:2000 --workers 3 --worker-class gevent
directory=/var/www/forum/server
autostart=true
autorestart=true
